# deploy/upstart.conf
# /etc/init/nag.conf

description "Nag tool"
author      "Tuomas Starck <tuomas@meetin.gs>"

start on (local-filesystems and net-device-up IFACE!=lo)
stop  on shutdown

respawn
respawn limit 5 5

setuid nobody
setgid nogroup

script
  set -e
  FIFO=`mktemp --dry-run`
  mkfifo $FIFO
  (logger -t nag -p local7.info <$FIFO &)
  exec >$FIFO
  rm $FIFO
  exec /usr/local/bin/nag 2>&1
end script
